# systemd/coilshield.service
# CoilShield ICCP controller — canonical systemd unit (flash-and-go install).
#
# Installed by deploy/bench-install.sh into /etc/systemd/system/coilshield.service.
# All paths and the runtime user come from the bench installer; do not hand-edit
# this file on the Pi. To change runtime config (Supabase URL, anon key, sim mode),
# edit /etc/coilshield/env and `sudo systemctl restart coilshield`.

[Unit]
Description=CoilShield ICCP controller (firmware)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=coilshield
Group=coilshield
WorkingDirectory=/opt/coilshield
EnvironmentFile=/etc/coilshield/env
ExecStart=/opt/coilshield/.venv/bin/iccp start
Restart=on-failure
RestartSec=10

# Hardware buses + cloud queue + serial cache live under these paths; firmware must
# not need write access elsewhere on the system root.
ReadWritePaths=/var/lib/coilshield /etc/coilshield
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
