# deploy/dashboard.service
# CoilShield ICCP — example systemd unit for the Flask dashboard (EXAMPLE)
#
# INSTALL (on the Pi):
#   1. sudo cp deploy/dashboard.service /etc/systemd/system/iccp-dashboard.service
#   2. Edit User, Group, WorkingDirectory, ExecStart, and Environment to match your install.
#      If you omit COILSHIELD_LOG_DIR here, ``iccp dashboard`` on Linux copies the controller's
#      telemetry dir from a running ``iccp start`` (``/proc``). Setting Environment= explicitly
#      still overrides that and is fine for non-default layouts.
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now iccp-dashboard
#
# Example (user onenous @ 0neNous-pi):
#   WorkingDirectory=/home/onenous/coilshield
#   ExecStart=/home/onenous/coilshield/.venv/bin/iccp dashboard --host 0.0.0.0 --port 8080
#   Environment=COILSHIELD_LOG_DIR=/home/onenous/coilshield/logs
#
# Note: the CLI auto-syncs systemd on recognized commands. For this dashboard unit,
# set a distinct ICCP_SYSTEMD_UNIT so the CLI reconciles the right service when you
# run `iccp dashboard` interactively elsewhere (e.g. Environment=ICCP_SYSTEMD_UNIT=iccp-dashboard),
# or set ICCP_SYSTEMD_SYNC=0 to skip systemctl from the CLI entirely.

[Unit]
Description=CoilShield ICCP web dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=onenous
Group=onenous
WorkingDirectory=/home/onenous/coilshield
ExecStart=/home/onenous/coilshield/.venv/bin/iccp dashboard --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
# Environment=COILSHIELD_LOG_DIR=/home/onenous/coilshield/logs

[Install]
WantedBy=multi-user.target
