# deploy/iccp.service
# CoilShield ICCP — systemd unit (compatibility path)
#
# Canonical example (edit there first): systemd/coilshield.service
# This copy exists so older docs and scripts that reference deploy/iccp.service
# still point at a valid unit. Keep [Unit]/[Service]/[Install] identical to systemd/coilshield.service.
#
# INSTALL (on the Pi):
#   sudo cp deploy/iccp.service /etc/systemd/system/iccp.service
#   # preferred canonical path:
#   sudo cp systemd/coilshield.service /etc/systemd/system/iccp.service
#
# Then edit User, Group, WorkingDirectory, and ExecStart if they differ.
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now iccp

[Unit]
Description=CoilShield ICCP controller (firmware)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=onenous
Group=onenous
WorkingDirectory=/home/onenous/coilshield
ExecStart=/home/onenous/coilshield/.venv/bin/iccp start
Restart=on-failure
RestartSec=10
Environment=COILSHIELD_SIM=0
# Environment=COILSHIELD_LOG_DIR=/home/onenous/coilshield/logs

[Install]
WantedBy=multi-user.target
